# monitor-test.service
[Unit]
Description=Monitor test process
After=network.target

[Service]
ExecStart=/usr/local/bin/monitor_test.sh
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target